# systemd/arkbot.service
[Unit]
Description=Arkanosis' Wikipedia bot

[Service]
Type=simple
ExecStart=arkbot

[Install]
WantedBy=default.target
