# systemd/arkbot.timer
[Unit]
Description=Hourly start of Arkanosis' Wikipedia bot

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
