# systemd/arkbot.service
[Unit]
Description=Arkbot

[Service]
Type=simple
ExecStart=arkbot

[Install]
WantedBy=default.target
